# furscan-tgbot.service
[Unit]
Description=FurScan Telegram Bot
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Desktop/pursuit
ExecStart=/home/user/Desktop/pursuit/.venv/bin/python tgbot.py
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/user/Desktop/pursuit/.env

[Install]
WantedBy=multi-user.target
